# 10-systemd/03-echod-lazy/echod.socket
[Unit]
Description=Echo daemon activation socket

[Socket]
ListenStream=/run/isel/echo-socket

[Install]
WantedBy=sockets.target

