# 10-systemd/01-service/useless.service
[Unit]
Description=Useless Service

[Service]
ExecStart=/home/isel/useless.sh
Restart=always
User=isel
Group=isel
